# root/etc/crontabs/root
*/15 * * * * s6-setuidgid abc python3 /usr/bin/feed2toot -l 5 -c /config/feed2toot.ini 2>&1